# services.d/memcached/run
#!/command/execlineb -P
memcached -u nobody