# services.d/iipsrv/run
#!/command/execlineb -P
with-contenv
/iipsrv/iipsrv.fcgi --bind 0.0.0.0:9000