# services.d/caddy/run
#!/command/execlineb -P
caddy run --config /etc/caddy/Caddyfile